# supervisord.conf
[supervisord]
nodaemon=true

[program:streamlit]
command=streamlit run frontend.py --server.port=8501 --server.address=0.0.0.0
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:flask]
command=gunicorn app:app --bind 0.0.0.0:5000
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
